# neomutt/.config/systemd/user/neomutt-clean-tmp.timer
[Unit]
Description=Run NeoMutt temp cleanup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
