# river/.config/systemd/user/river-session.target
[Unit]
Description=User services for river
BindsTo=graphical-session.target
