# systemd-timer/.config/systemd/user/backup-wow-wtf.service
[Unit]
Description=WoW WTF (Settings) Folder Backup

[Service]
Type=oneshot
ExecStart=/home/shmcg/.local/bin/backup-wow-wtf.sh
