# dropbox/.config/systemd/user/dropbox.service.d/override.conf
[Service]
Environment=WAYLAND_DISPLAY=wayland-1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority

[Unit]
After=waybar@LG_ULTRA.service
Wants=waybar@LG_ULTRA.service
